# data/shoots.cron
0 0 1 1 1-5 veveshooter-now >/dev/null 2>&1
0 0 1 1 1-5 veveshooter-now >/dev/null 2>&10 0 1 1 1-5 veveshooter-now >/dev/null 2>&1
0 0 1 1 1-5 veveshooter-now >/dev/null 2>&1
0 0 1 1 1-5 veveshooter-now >/dev/null 2>&1
0 0 1 1 1-5 veveshooter-now >/dev/null 2>&1
